# etc/supervisord.conf
[program:czcake8080]
command=python /home/www/czcake/src/manager.py
autostart=true
stderr_logfile = /var/log/nginx/app.log
stdout_logfile = /var/log/nginx/app.log